# supervisord-confs/app.conf
[program:app]
command = bash /app.sh
user = root
autostart = true
autorestart = true
stderr_logfile = /var/log/supervisor/app-stderr.log
stdout_logfile = /var/log/supervisor/app-stdout.log
